# deploy/systemd/tgaggerator-api.service
[Unit]
Description=TGaggregator API service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgagg
Group=tgagg
WorkingDirectory=/opt/tgaggerator
EnvironmentFile=/opt/tgaggerator/.env
ExecStart=/usr/bin/env uv run python scripts/run_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
