# deploy/systemd/tgaggerator-collector.service
[Unit]
Description=TGaggregator collector service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tgagg
Group=tgagg
WorkingDirectory=/opt/tgaggerator
EnvironmentFile=/opt/tgaggerator/.env
ExecStart=/usr/bin/env uv run python -m tgaggerator.cli ingest-loop --interval 30
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
